# famick/rootfs/etc/services.d/scheduler/run
#!/command/with-contenv bash
# s6 service: supercronic scheduler. Same binary the docker-compose
# scheduler sidecar runs, fed by the same crontab the canonical image
# bakes in at /app/scheduler-crontab. supercronic invokes the web app
# in `run-job <name>` mode for each cron tick.

set -euo pipefail

until pg_isready -h localhost -U famick -d famick >/dev/null 2>&1; do
    sleep 1
done

TENANT_ID=$(cat /data/config/tenant-id)

export ASPNETCORE_ENVIRONMENT=Production
export ConnectionStrings__DefaultConnection="Host=localhost;Database=famick;Username=famick;Password=famick"
export FixedTenantId="${TENANT_ID}"
export SelfHosted__TenantId="${TENANT_ID}"
export SelfHosted__ApplicationName="Famick (Home Assistant)"
export JwtSettings__RsaPrivateKeyPemFile=/data/keys/jwt-rsa.pem
export Storage__Path=/data

cd /app
exec supercronic /app/scheduler-crontab
